# node-app.service
[Unit]
Description=Node testing

[Service]
ExecStart=/usr/bin/node /vagrant/app.js

[Install]
WantedBy=multi-user.target

